# configs/flask_celery.conf
description "celery to run flask_google background tasks"

start on runlevel [2345]
stop on runlevel [!2345]

setuid vagrant
setgid www-data

env PATH=/home/vagrant/flask_google/flask_env/bin
chdir /home/vagrant/flask_google/
exec celery  worker -A app.celery --loglevel=debug
